# systemd/frontend.service
[Unit]
Description=Student Management System - Next.js Frontend
After=network.target

[Service]
Type=simple
User=app
Group=app
WorkingDirectory=/home/app/student-management/frontend
Environment="NODE_ENV=production"
Environment="PORT=3000"
ExecStart=/usr/bin/node /home/app/student-management/frontend/node_modules/.bin/next start -p 3000
Restart=on-failure
RestartSec=5
StandardOutput=append:/home/app/student-management/logs/frontend-stdout.log
StandardError=append:/home/app/student-management/logs/frontend-stderr.log

[Install]
WantedBy=multi-user.target
